# udev-rules/99-auto-cpu-governor-powersave.rules
# AC UNPLUGGED
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/usr/local/bin/cpu-governor-powersave"
